# systemd-evdev-proxy/evdev-proxy.service
[Unit]
Description=Run evdev-proxy.

[Service]
Type=simple
ExecStart=/usr/bin/evdev-proxy
